# root/etc/services.d/tracker/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc /bin/bash /config/tracker-add-auto.sh
