# root/etc/services.d/transmission/run
#!/usr/bin/with-contenv bash

s6-setuidgid abc /usr/bin/transmission-daemon -g /config -c /watch -f
