# root/etc/s6-overlay/s6-rc.d/init-transmission-config/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash


# copy config
if [[ ! -f /config/settings.json ]]; then
    cp /defaults/settings.json /config/settings.json
fi

# make folders, optionally set paths in config
mkdir -p /watch

INIT_COMPLETE_PATH="/downloads/complete"
INIT_INCOMPLETE_PATH="/downloads/incomplete"

if [[ -n "$COMPLETE_PATH" ]]; then
    INIT_COMPLETE_PATH="$COMPLETE_PATH"
    sed -i "/\"download-dir\"/c\    \"download-dir\": \"$COMPLETE_PATH\"," /config/settings.json
fi
mkdir -p "$INIT_COMPLETE_PATH"

if [[ -n "$INCOMPLETE_PATH" ]]; then
    INIT_INCOMPLETE_PATH="$INCOMPLETE_PATH"
    sed -i "/\"incomplete-dir\"/c\    \"incomplete-dir\": \"$INCOMPLETE_PATH\"," /config/settings.json
fi
mkdir -p "$INIT_INCOMPLETE_PATH"

if [[ -n "$USER" ]] && [[ -n "$PASS" ]]; then
    sed -i '/rpc-authentication-required/c\    "rpc-authentication-required": true,' /config/settings.json
    sed -i "/rpc-username/c\    \"rpc-username\": \"$USER\"," /config/settings.json
    sed -i "/rpc-password/c\    \"rpc-password\": \"$PASS\"," /config/settings.json
else
    sed -i '/rpc-authentication-required/c\    "rpc-authentication-required": false,' /config/settings.json
    sed -i "/rpc-username/c\    \"rpc-username\": \"$USER\"," /config/settings.json
    sed -i "/rpc-password/c\    \"rpc-password\": \"$PASS\"," /config/settings.json
fi

if [[ -n "$WHITELIST" ]]; then
    sed -i '/rpc-whitelist-enabled/c\    "rpc-whitelist-enabled": true,' /config/settings.json
    sed -i "/\"rpc-whitelist\"/c\    \"rpc-whitelist\": \"$WHITELIST\"," /config/settings.json
else
    sed -i '/rpc-whitelist-enabled/c\    "rpc-whitelist-enabled": false,' /config/settings.json
    sed -i "/\"rpc-whitelist\"/c\    \"rpc-whitelist\": \"$WHITELIST\"," /config/settings.json
fi

if [[ -n "$HOST_WHITELIST" ]]; then
    sed -i '/rpc-host-whitelist-enabled/c\    "rpc-host-whitelist-enabled": true,' /config/settings.json
    sed -i "/\"rpc-host-whitelist\"/c\    \"rpc-host-whitelist\": \"$HOST_WHITELIST\"," /config/settings.json
else
    sed -i '/rpc-host-whitelist-enabled/c\    "rpc-host-whitelist-enabled": false,' /config/settings.json
    sed -i "/\"rpc-host-whitelist\"/c\    \"rpc-host-whitelist\": \"$HOST_WHITELIST\"," /config/settings.json
fi

if [[ -n "${PEERPORT}" ]]; then
    sed -i "/\"peer-port\"/c\    \"peer-port\": ${PEERPORT}," /config/settings.json
    sed -i '/peer-port-random-on-start/c\     "peer-port-random-on-start": false,' /config/settings.json
fi

# Handle old theme locations
mkdir -p {/transmissionic,/combustion-release,/flood-for-transmission,/kettu,/transmission-web-control}
echo /transmissionic /combustion-release /flood-for-transmission /kettu /transmission-web-control | xargs -n1 ln -s /defaults/index.html

# permissions
lsiown abc:abc \
    /config/settings.json

if [[ "$(stat -c '%U' /downloads)" != "abc" ]]; then
    lsiown abc:abc /downloads
fi

if [[ "$(stat -c '%U' "$INIT_COMPLETE_PATH")" != "abc" ]]; then
    lsiown abc:abc "$INIT_COMPLETE_PATH"
fi

if [[ "$(stat -c '%U' "$INIT_INCOMPLETE_PATH")" != "abc" ]]; then
    lsiown abc:abc "$INIT_INCOMPLETE_PATH"
fi

if [[ "$(stat -c '%U' /watch)" != "abc" ]]; then
    lsiown abc:abc /watch
fi
