# root/etc/s6-overlay/s6-rc.d/init-transmission-config/up
/etc/s6-overlay/s6-rc.d/init-transmission-config/run
