# root/etc/s6-overlay/s6-rc.d/svc-transmission/finish
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

pid=$(pidof transmission-daemon)

/bin/kill -s STOP ${pid}

# tail --pid=${pid} -f /dev/null
## This is causing an s6 timeout error.
