# root/etc/s6-overlay/s6-rc.d/svc-cron/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

exec busybox crond -f -S -l 5
